# config/news-sentry-realtime.service
# News Sentry — 实时雷达单次运行服务（被 timer 触发）
# 不要手动启动此服务，由 news-sentry-realtime.timer 触发

[Unit]
Description=News Sentry — 实时雷达单次采集周期

[Service]
Type=oneshot
User=news-sentry
Group=news-sentry
WorkingDirectory=/opt/news-sentry

# 全链路 bounded run（单 target）
ExecStart=/opt/news-sentry/.venv/bin/python3 -m news_sentry.cli run --target italy --stage all

# 环境变量
EnvironmentFile=-/opt/news-sentry/.env

# 超时保护（单次运行不超过 5 分钟）
TimeoutStartSec=300

# 安全加固
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/news-sentry/data
