# config/news-sentry-realtime.service
# News Sentry realtime collection service.
# Triggered by news-sentry-realtime.timer; do not start manually during normal ops.

[Unit]
Description=News Sentry realtime collection cycle
After=network-online.target news-sentry.service
Wants=network-online.target

[Service]
Type=oneshot
User=newssentry
Group=newssentry
WorkingDirectory=/opt/news-sentry/production/repo
EnvironmentFile=/opt/news-sentry/production/.env
ExecStart=/opt/news-sentry/production/repo/tools/run_realtime_collection.sh

# One full multi-target cycle should normally finish well under this.
TimeoutStartSec=900

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/srv/news-sentry/production /opt/news-sentry/production
