# scripts/news-sentry.service
# News Sentry — systemd user service unit
#
# 安装（用户级，无需 sudo）：
#   mkdir -p ~/.config/systemd/user
#   cp scripts/news-sentry.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now news-sentry.service
#
# 常用命令：
#   systemctl --user status news-sentry   # 查看状态
#   systemctl --user restart news-sentry  # 重启服务
#   journalctl --user -u news-sentry -f   # 查看日志

[Unit]
Description=News Sentry — News Intelligence Monitor
Documentation=https://github.com/xucroyuri/news-sentry
After=network.target

[Service]
Type=simple
ExecStart=%h/.news-sentry/venv/bin/python -m news_sentry.cli serve --foreground
WorkingDirectory=%h/.news-sentry
Restart=always
RestartSec=10

# 日志输出（追加写入，由 systemd 管理轮转）
StandardOutput=append:%h/.news-sentry/logs/serve.log
StandardError=append:%h/.news-sentry/logs/serve.log

# 环境变量默认值（可在 ~/.news-sentry/.env 中覆盖）
Environment=NEWSSENTRY_DATA_DIR=%h/.news-sentry/data
Environment=NEWSSENTRY_LOG_DIR=%h/.news-sentry/logs

[Install]
WantedBy=default.target
